# etc/crontab
@hourly bash -l -c 'backup &> ~/.backup-log'
@daily bash -l -c 'backup clean &> ~/.backup-clean-log'
@daily bash -l -c 'backup purge'
