# etc/crontab
@hourly backup
@daily backup clean
@daily backup purge
